# YoGoKoBox/etc/systemd/system/its-services.target.wants/its-services-pre.target
[Unit]
Description=YoGoKo ITS Services synchronization point
Before=its-services.target

[Install]
WantedBy=its-services.target
