# YoGoKoBox/etc/systemd/system/multi-user.target.wants/its-services.target
[Unit]
Description=YoGoKo ITS Services
Requires=its-services-pre.target
After=its-services-pre.target

[Install]
WantedBy=multi-user.target
